# autostart/ioturing.sh.desktop
[Desktop Entry]
Exec=IoTuring
Icon=dialog-scripts
Name=IoTuring
Type=Application
X-KDE-AutostartScript=true
